# minidsp/rootfs/etc/services.d/minidsp/run
#!/usr/bin/env bash
set -e

echo "[miniDSP] Starting Gunicorn (PID 1 будет именно он)..."

source /opt/venv/bin/activate

# ←←←← ВОТ ЭТО ГЛАВНОЕ: exec + s6-setuidgid !!!
exec s6-setuidgid abc:abc gunicorn --bind 0.0.0.0:8080 \
     --workers 2 \
     --access-logfile - \
     --error-logfile - \
     api:app